# scripts/starnion-ui.service
[Unit]
Description=Starnion UI (Next.js)
Documentation=https://github.com/jikime/starnion
# UI는 게이트웨이 API를 호출하므로 게이트웨이가 먼저 올라온 후 시작.
After=network-online.target starnion-gateway.service
Wants=network-online.target
PartOf=starnion.target

[Service]
Type=simple
User=root
Group=root

# ── 실행 경로 ────────────────────────────────────────────────────────────────
# standalone 빌드: server.js가 WorkingDirectory 루트에 위치
# (starnion update가 자동으로 node 경로를 조정함)
WorkingDirectory=/root/.starnion/ui
ExecStart=/usr/bin/node server.js

# ── 환경 변수 ────────────────────────────────────────────────────────────────
Environment=HOME=/root
Environment=NODE_ENV=production
Environment=PORT=3893
Environment=HOSTNAME=0.0.0.0

# Next.js / NextAuth 환경 변수는 빌드 시 .env 에서 번들됨.
# 런타임에 바꿔야 할 값만 여기에 추가.
# AUTH_TRUST_HOST=true 는 nginx 프록시 뒤에서 필수 (ui/.env 에서 설정).
# EnvironmentFile=/etc/starnion/ui.env

# ── 재시작 정책 ──────────────────────────────────────────────────────────────
Restart=always
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=30

# ── 리소스 제한 ──────────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096

# ── 종료 ─────────────────────────────────────────────────────────────────────
TimeoutStopSec=15
KillMode=control-group
KillSignal=SIGTERM

# ── 보안 ─────────────────────────────────────────────────────────────────────
NoNewPrivileges=true

# ── 로그 ─────────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=starnion-ui

[Install]
WantedBy=starnion.target multi-user.target
